# s6/proftpd/run
#!/bin/sh

[ -d /run/proftpd ] || mkdir -pv /run/proftpd
exec /usr/local/sbin/proftpd --nodaemon --config /usr/local/etc/proftpd.conf
