# s6/proftpd/run
#!/bin/sh

[ -d /run/proftpd ] || mkdir -pv /run/proftpd
exec /usr/sbin/proftpd --ipv6 --nodaemon --config /etc/proftpd/proftpd.conf
